# packaging/pfsense/files/usr/local/etc/rc.d/zid_packages
#!/bin/sh

# PROVIDE: zid_packages
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="zid_packages"
rcvar="zid_packages_enable"

command="/usr/local/sbin/zid-packages"
command_args="daemon"
pidfile="/var/run/zid-packages.pid"

: ${zid_packages_enable:="NO"}

start_cmd="zid_packages_start"
stop_cmd="zid_packages_stop"

zid_packages_start() {
	if [ "${zid_packages_enable}" != "YES" ]; then
		echo "zid_packages disabled in rc.conf"
		return 0
	fi
	/usr/sbin/daemon -f -p "${pidfile}" ${command} ${command_args}
}

zid_packages_stop() {
	if [ -f "${pidfile}" ]; then
		PID=$(cat "${pidfile}" 2>/dev/null || true)
		if [ -n "${PID}" ]; then
			kill "${PID}" 2>/dev/null || true
		fi
		rm -f "${pidfile}" 2>/dev/null || true
	fi
	pkill -f "${command} ${command_args}" 2>/dev/null || true
}

run_rc_command "$1"
